# src/msg.conf
[program:msg]
command=php /PATH/TO/msg/bin/receive.php
user=USER
stderr_logfile=/var/log/supervisor/msg.err.log
stdout_logfile=/var/log/supervisor/msg.out.log
autostart=true
autorestart=true
